# scanner/scan_whois.cron
#! /bin/sh
* 23 * * * /user/bin/python3 start.py
